# homepage/rootfs/etc/services.d/homepage/finish
#!/usr/bin/with-contenv bashio

bashio::log.info "Homepage service stopped."

exit 0
