# homepage/rootfs/etc/services.d/homepage/run
#!/usr/bin/with-contenv bashio

# Read options
ALLOWED_HOSTS=$(bashio::config 'allowed_hosts')

# Set environment variables
export HOMEPAGE_ALLOWED_HOSTS="${ALLOWED_HOSTS}"

# Log startup
bashio::log.info "Starting Homepage with allowed hosts: ${ALLOWED_HOSTS}"

if [ -d "/app/config" ]; then
    rm -rf /app/config
    bashio::log.info "Delete existing /app/config"
fi

ln -s /config /app/config

cd /app

# Start Homepage app
pnpm start
